# supervisord.conf
[supervisord]
[program:playground]
command=./odin-playground.bin
stopsignal=INT
stdout_logfile=./odin-playground.stdout.log
stderr_logfile=./odin-playground.stderr.log
