# supervisord.conf
[supervisord]
[program:playground]
command=./odin-playground.bin
