# supervisord.conf
[supervisord]
[program:playground]
command=./odin-playground.bin
stopsignal=INT
stdout_logfile=./odin-playground.stdout.log
redirect_stderr=true
